# packaging/deb/prerm
#!/bin/sh
# prerm for send-to-influx.deb
set -e

# prerm is only ever invoked as remove/upgrade/deconfigure/failed-upgrade -
# there is no prerm "purge" action (purge-time cleanup lives in postrm).
if [ -d /run/systemd/system ] && [ "$1" = "remove" ]; then
    systemctl stop send-to-influx || true
    systemctl disable send-to-influx || true
fi

exit 0
